# dot_config/systemd/user/xray-mitm.service
[Unit]
Description=Xray MITM Domain Fronting Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Note: '~' does not work in systemd ExecStart. We use '%h' which resolves to the user's home directory.
ExecStart=/usr/bin/xray run -c %h/.config/xray/MITM-DomainFronting.json

# Restart the service if it crashes
Restart=on-failure
RestartSec=5s

# Optional: Limit logs to prevent filling up your disk if it gets stuck in a restart loop
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
